# roles/archlinux/templates/reflector.timer
[Unit]
Description=Pacman mirrorlist update

[Timer]
OnCalendar=*-*-* 00:00:00

[Install]
WantedBy=default.target
